# mavaia-api.service
[Unit]
Description=Mavaia API Server Daemon
After=network.target

[Service]
Type=simple
User=mike
Group=mike
WorkingDirectory=/home/mike/Mavaia
Environment=PYTHONPATH=/home/mike/Mavaia
Environment=MAVAIA_ENABLE_HEAVY_MODULES=true
ExecStart=/home/mike/Mavaia/.venv/bin/python3 /home/mike/Mavaia/scripts/mavaia_api_daemon.py
Restart=always
RestartSec=10
StandardOutput=append:/home/mike/Mavaia/api_daemon.log
StandardError=append:/home/mike/Mavaia/api_daemon.log

[Install]
WantedBy=multi-user.target
