# oricli-api.service
[Unit]
Description=Oricli-Alpha API Server Daemon
After=network.target

[Service]
Type=simple
User=mike
Group=mike
WorkingDirectory=/home/mike/Oricli-Alpha
Environment=PYTHONPATH=/home/mike/Oricli-Alpha
Environment=MAVAIA_ENABLE_HEAVY_MODULES=true
Environment=OLLAMA_URL=http://localhost:11434
Environment=OLLAMA_MODEL=qwen2.5:7b
ExecStart=/home/mike/Oricli-Alpha/.venv/bin/python3 /home/mike/Oricli-Alpha/scripts/oricli_api_daemon.py
Restart=always
RestartSec=10
StandardOutput=append:/home/mike/Oricli-Alpha/api_daemon.log
StandardError=append:/home/mike/Oricli-Alpha/api_daemon.log

[Install]
WantedBy=multi-user.target
