# oricli-backbone.service
[Unit]
Description=Oricli-Alpha Go Backbone (Hive Orchestrator)
After=network.target

[Service]
Type=simple
User=mike
Group=mike
WorkingDirectory=/home/mike/Mavaia
ExecStart=/home/mike/Mavaia/scripts/start_go_hive.sh
Restart=always
RestartSec=15
LimitNOFILE=65535
StandardOutput=append:/home/mike/Mavaia/go_backbone.log
StandardError=append:/home/mike/Mavaia/go_backbone.log

[Install]
WantedBy=multi-user.target
