# oricli-backbone.service
[Unit]
Description=Oricli-Alpha Go Backbone (Hive Orchestrator)
After=network.target

[Service]
Environment=OLLAMA_URL=http://localhost:11434
Environment=OLLAMA_MODEL=qwen2.5-coder:3b
Environment=ORICLI_BROWSERLESS_URL=ws://localhost:3000
Environment=MAVAIA_S3_ENDPOINT=http://localhost:9000
Environment=AWS_ACCESS_KEY_ID=oricli-admin
Environment=AWS_SECRET_ACCESS_KEY=oricli-sovereign-2025
Environment=AWS_BUCKET_NAME=oricli-state
Environment=AWS_DEFAULT_REGION=us-east-1
Type=simple
User=mike
Group=mike
WorkingDirectory=/home/mike/Mavaia
ExecStart=/home/mike/Mavaia/scripts/start_go_hive_pure.sh
Restart=always
RestartSec=15
LimitNOFILE=65535
StandardOutput=append:/home/mike/Mavaia/go_backbone.log
StandardError=append:/home/mike/Mavaia/go_backbone.log

[Install]
WantedBy=multi-user.target
