# docker/browserless/oricli-browserless.service
[Unit]
Description=Oricli Browserless Chrome CDP Service
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/mike/Mavaia/docker/browserless
ExecStart=/usr/bin/docker compose up -d --pull missing
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=180
TimeoutStopSec=30
User=mike

[Install]
WantedBy=multi-user.target
